# media-check/crontab
0 3 * * * cd /app && python check.py 2>&1
