# sonarr-import-resolver/crontab
*/15 * * * * cd /app && python resolve.py 2>&1
